# lscustoms.service
[Unit]
Description=LSCustoms Services (Astro + Node.js)
After=network.target

[Service]
Type=forking
User=ec2-user
WorkingDirectory=/home/ec2-user/Proyecto LSCustoms
ExecStart=/home/ec2-user/Proyecto LSCustoms/start-services.sh

# Restart configuration
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target 